# discord-bot/discord-bot.service.template
[Unit]
Description=Minecraft Discord Bot
After=network.target

[Service]
User=__USER__
WorkingDirectory=__INSTALL_DIR__/discord-bot
EnvironmentFile=__INSTALL_DIR__/discord-bot/.env
ExecStart=/usr/bin/python3 bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
